# cookbooks/cla_clamav/templates/default/clamav.cron.erb
#!/bin/sh

rm -rf /tmp/clamav*

/usr/bin/freshclam --quiet --log=/tmp/clamav.log
/usr/bin/clamscan --scan-mail=no --no-summary --infected --log=/tmp/clamav.log --recursive <%= node[:clamav][:paths].join(" ") %> 2>/dev/null

if [ "$?" -ne "0" ]; then
  cat - /tmp/clamav.log << EOF | /usr/sbin/sendmail <%= node[:clamav][:email_to] %>
To: <%= node[:clamav][:email_to] %>
From: nobody@`/bin/hostname -f`
Subject: Virus detected on `/bin/hostname -f`

EOF

fi
