# cookbooks/cla_shell/files/default/profile.d/maven-2-0-9.sh
#!/bin/sh

if [ -d /opt/maven/current/bin ]; then
 M2_HOME=/opt/maven/current
 M2=${M2_HOME}/bin
 export M2 M2_HOME
 PATH=${M2}:${PATH}
 export PATH
fi
